# config/supervisord.conf
[supervisord]
#nodaemon=true

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:80        ; (ip_address:port specifier, *:port for all iface)
;username=daedalus              ; (default is no username (open server))
;password=daedalus               ; (default is no password (open server))

[program:event_data_logger]
command=/usr/bin/python3 datalogger.py --data /usr/local/daedalus/data/unified --backup /mnt/data/unified
directory=/usr/local/daedalus/code
autostart=false
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=1

[program:event_based_camera]
command= /usr/bin/python3 event_camera_datalogger.py --data /usr/local/daedalus/data/pi_events_data --backup /mnt/data/pi_events_data --socket /tmp/event.sock 00051501 00051503 00051505 00051500
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=1

[program:imu_servo]
command= /usr/bin/python3 servo_imu.py
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=1

[program:i_m_u]
command= /usr/bin/python3 imu_datalogger.py --data /usr/local/daedalus/data/imu --backup /mnt/data/imu --socket /tmp/imu.sock --i2c_address 0x69
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=3

[program:pi_picture_camera]
command= /usr/bin/python3 pi_camera_picture_datalogger.py --data /usr/local/daedalus/data/pi_picture --backup /mnt/data/pi_picture --camera 0
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=2

[program:pi_video_camera]
command= /usr/bin/python3 pi_camera_video_datalogger.py --data /usr/local/daedalus/data/pi_video --backup /mnt/data/pi_video --camera 1
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=2

[program:g_p_s]
command= /usr/bin/python3 gps_datalogger.py --data /usr/local/daedalus/data/gps_data --backup /mnt/data/gps_data --socket /tmp/gps.sock --port /dev/ttyACM0
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=3

[program:atmos_temp_sensor]
command= /usr/bin/python3 atmos_datalogger.py --data /usr/local/daedalus/data/atmos_data --backup /mnt/data/atmos_data --socket /tmp/atmos.sock --i2c_address 0x77
directory=/usr/local/daedalus/code
autorestart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=3

[program:o_l_e_d]
command= /usr/bin/python3 data_monitor.py
directory=/usr/local/daedalus/code
autorestart=true
autostart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=2

[program:r_f_transmitter]
command= /usr/bin/python3 rf_transmitter.py --port /dev/ttyUSB0
directory=/usr/local/daedalus/code
autorestart=true
autostart=true
#startsecs=0
startretries=10000
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=2

[program:event_camera_tester]
command=/usr/bin/python3 event_camera_streamer.py --port 8000 --serial 00050869
directory=/usr/local/daedalus/code
autostart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=5

[program:pi_camera_tester]
command=/usr/bin/python3 pi_camera_streamer.py --port 8001 --camera 0
directory=/usr/local/daedalus/code
autostart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=5

[program:ir_camera_tester]
command=/usr/bin/python3 ir_camera_streamer.py --port 8002
directory=/usr/local/daedalus/code
autostart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
priority=5
